# root/etc/services.d/peoplesoft/run
#!/usr/bin/with-contenv sh

#exec s6-setuidgid chrome echo $NPM_CONFIG_PREFIX
exec s6-setuidgid chrome /usr/src/app/src/index.ts